# nomouse.service
[Unit]
Description=Nomouse - Transparent Grid Overlay
After=network.target display-manager.service
Wants=display-manager.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/ishan/nomouse
ExecStart=/usr/bin/python3 /home/ishan/nomouse/nomouse.py
Restart=always
RestartSec=3
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/ishan/.Xauthority

[Install]
WantedBy=multi-user.target 